# deploy/h59.oracle.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>h59.oracle.morning</string>
    <!-- 매일 7:00 — 베르 조간 발행물 합성. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-s</string>
        <string>-X</string>
        <string>POST</string>
        <string>--max-time</string>
        <string>300</string>
        <string>http://localhost:8001/briefing/run?kind=morning</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/chocolat/data/logs/oracle-morning.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chocolat/data/logs/oracle-morning.log</string>
</dict>
</plist>
