# deploy/h59.oracle.mongo-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>h59.oracle.mongo-backup</string>

    <!-- 매주 일요일 03:00 (Weekday: 0=일) — Oracle MongoDB 주간 백업 -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/chocolat/projects/oracle/deploy/backup-mongo.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/chocolat/data/logs/oracle-mongo-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chocolat/data/logs/oracle-mongo-backup.log</string>
</dict>
</plist>
